# openwrt/luci-app-shawnwrt-quickstart/root/etc/uci-defaults/50_shawnwrt-quickstart
#!/bin/sh

uci -q set luci.main.default_url='/cgi-bin/luci/admin/quickstart'
uci -q commit luci

rm -f /tmp/luci-indexcache /tmp/luci-indexcache.*
exit 0
