# openwrt/luci-app-shawnwrt-index/root/etc/uci-defaults/50_shawnwrt-quickstart
#!/bin/sh

uci -q set luci.main.default_url='/cgi-bin/luci/admin/index'
uci -q commit luci

rm -f /tmp/luci-indexcache /tmp/luci-indexcache.* /tmp/luci-modulecache/*
exit 0
